# project/configs/smm_bicycle_celery.service
[Unit]
Description=SMM Bicycle Celery

[Service]
ExecStart=/root/projects/project/.env/bin/celery -A tasks worker
